# scripts/init.d/vserial
#!/sbin/openrc-run

depend() {
   need localmount
}

: ${VDEVICE:="/dev/ttyS0"}
: ${VNUMBER:="2"}
: ${VBAUDRATE:="9600"}
: ${VLINE:="no"}

command="/usr/bin/vserial"

case "${VLINE}" in
    [Yy][Ee][Ss]|[Tt][Rr][Uu][Ee]|1)
        vline_arg="-line"
        ;;
    *)
        vline_arg=""
        ;;
esac

command_args="-baudrate ${VBAUDRATE} -device ${VDEVICE} -number ${VNUMBER} ${vline_arg}"
command_background="yes"
pidfile="/run/vserial.pid"
